# ganglia/templates/gmetad.init.erb
#!/bin/sh
# $Id: gmetad.init 156 2002-10-18 21:57:58Z sacerdoti $
#
# chkconfig: 2345 20 80
# description: gmetad startup script
#
GMETAD=/usr/sbin/gmetad_<%= gmetad_name %>
. /etc/rc.d/init.d/functions

RETVAL=0

case "$1" in
   start)
      echo -n "Starting GANGLIA gmetad_<%= gmetad_name %>: "
      [ -f $GMETAD ] || exit 1
      
       mkdir /var/lib/ganglia/<%= gmetad_name %>/rrds >& /dev/null && chown nobody:nobody /var/lib/ganglia/<%= gmetad_name %>/rrds
      daemon $GMETAD -c /etc/ganglia/gmetad_<%= gmetad_name %>.conf
      RETVAL=$?
      echo
      [ $RETVAL -eq 0 ] && touch /var/lock/subsys/gmetad_<%= gmetad_name %>
      ;;

  stop)
      echo -n "Shutting down GANGLIA gmetad_<%= gmetad_name %>: "
      killproc gmetad_<%= gmetad_name %>
      RETVAL=$?
      echo
      [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/gmetad_<%= gmetad_name %>
      ;;

  restart|reload)
      $0 stop
      $0 start
      RETVAL=$?
      ;;
  status)
      status gmetad_<%= gmetad_name %>
      RETVAL=$?
      ;;
  *)
      echo "Usage: $0 {start|stop|restart|status}"
      exit 1
esac

exit $RETVAL
